# crontab
15 * * * * root /sps30_collect_data.py
3 18 * * 6 root /sps30_sensor_clean.py
#