# app/crontab
# this is an example file and the cd /path/... needs to be change to match your local setup
#USING PYTHON VENV: Don't activate the env, but use the python installed in it...
#/home/ubuntu/scripts/venv/bin/python3
# add to sudo crontab -e

* * * * * cd /home/pi/IoT-Microservice/app/oracle/ && ./sps30_collect_data_oneshot.py
30 4 * * SUN cd /home/pi/git-repos/IoT-Microservice/app/oracle/ && ./sps30_sensor_clean_oneshot.py

* * * * * cd /home/pi/git-repos/IoT-Microservice/app/oracle/ && ./scd30_collect_data.py

* * * * * cd /home/pi/git-repos/IoT-Microservice/app/oracle/ && ./scd41_collect_data.py

# blank line needed att end of file to function properly
